# ansible-role-i2b2/templates/etc/profile.d/ant.sh
# {{ ansible_managed }}
export ANT_HOME=/opt/{{ i2b2_artifacts.ant.dir }}
export PATH=$ANT_HOME/bin:$PATH
