# ansible-role-i2b2/templates/etc/profile.d/jboss.sh
# {{ ansible_managed }}
export JBOSS_HOME=/opt/{{ i2b2_artifacts.jboss.dir }}
